# settings/linux/crontab/idachev-vmw
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.WVYRRW/crontab installed on Fri Jun  7 13:19:15 2013)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
# Edit this file to introduce tasks to be run by cron.
# 
# Each task to run has to be defined through a single line
# indicating with different fields when the task will be run
# and what command to run for the task
# 
# To define the time you can provide concrete values for
# minute (m), hour (h), day of month (dom), month (mon),
# and day of week (dow) or use '*' in these fields (for 'any').# 
# Notice that tasks will be started based on the cron's system
# daemon's notion of time and timezones.
# 
# Output of the crontab jobs (including errors) is sent through
# email to the user the crontab file belongs to (unless redirected).
# 
# For example, you can run a backup of all your user accounts
# at 5 a.m every week with:
# 0 5 * * 1 tar -zcf /var/backups/home.tgz /home/
# 
# For more information see the manual pages of crontab(5) and cron(8)
# 
# m h  dom mon dow   command
MAILTO="idachev@vmware.com"

# Backups local
05 10 * * * . /home/idachev/.keychain/*-sh && . /home/idachev/.keychain/*-sh-gpg && /usr/bin/env python /home/idachev/bin/backup/backup_work.py -d debug --inc-rsync=/home/idachev/bin/inc-rsync --src-list=/home/idachev/bin/backup/backup_src_list_local.txt --dst-list=/home/idachev/bin/backup/backup_dst_list_local.txt >> /tmp/cron_job_backup_local.log 2>&1
05 14 * * * . /home/idachev/.keychain/*-sh && . /home/idachev/.keychain/*-sh-gpg && /usr/bin/env python /home/idachev/bin/backup/backup_work.py -d debug --inc-rsync=/home/idachev/bin/inc-rsync --src-list=/home/idachev/bin/backup/backup_src_list_local.txt --dst-list=/home/idachev/bin/backup/backup_dst_list_local.txt >> /tmp/cron_job_backup_local.log 2>&1
05 18 * * * . /home/idachev/.keychain/*-sh && . /home/idachev/.keychain/*-sh-gpg && /usr/bin/env python /home/idachev/bin/backup/backup_work.py -d debug --inc-rsync=/home/idachev/bin/inc-rsync --src-list=/home/idachev/bin/backup/backup_src_list_local.txt --dst-list=/home/idachev/bin/backup/backup_dst_list_local.txt >> /tmp/cron_job_backup_local.log 2>&1
05 23 * * * . /home/idachev/.keychain/*-sh && . /home/idachev/.keychain/*-sh-gpg && /usr/bin/env python /home/idachev/bin/backup/backup_work.py -d debug --inc-rsync=/home/idachev/bin/inc-rsync --src-list=/home/idachev/bin/backup/backup_src_list_local.txt --dst-list=/home/idachev/bin/backup/backup_dst_list_local.txt >> /tmp/cron_job_backup_local.log 2>&1

# Backup cloud
#30 02 * * * . /home/idachev/.keychain/*-sh && . /home/idachev/.keychain/*-sh-gpg && /usr/bin/env python /home/idachev/bin/backup/backup_work.py -d debug --inc-rsync=/home/idachev/bin/inc-rsync --src-list=/home/idachev/bin/backup/backup_src_list_cloud.txt --dst-list=/home/idachev/bin/backup/backup_dst_list_cloud.txt >> /tmp/cron_job_backup_cloud.log 2>&1

# Some web sites during the day
45 10 * * * env DISPLAY=:0 /home/idachev/bin/dbus_crontab_exec.sh chrome /opt/google/chrome/chrome 'http://antwrp.gsfc.nasa.gov/apod/' > /dev/null 2>&1
30 14 * * * env DISPLAY=:0 /home/idachev/bin/dbus_crontab_exec.sh chrome /opt/google/chrome/chrome 'http://www.kaldata.com/index.php?action=news&catid=2' > /dev/null 2>&1
30 16 * * * env DISPLAY=:0 /home/idachev/bin/dbus_crontab_exec.sh chrome /opt/google/chrome/chrome 'http://www.linuxtoday.com/' > /dev/null 2>&1

# Show my wallpapers every 5 minutes
#*/5 * * * *  /home/idachev/bin/dbus_crontab_exec.sh mate /home/idachev/bin/wallpaper.sh /home/idachev/Pictures > /dev/null 2>&1

# Update xfce4 wallpapers every 5 minutes
*/5 * * * * DISPLAY=:0.0 xfdesktop --reload > /dev/null 2>&1

#00 10 * * 4 echo "There is need each week to do an integration from vmkernal-main to op-licensing-vpxd branch." | mail -s "Weekly branch integration notification" licensing-dev@vmware.com > /dev/null 2>&1

01 * * * * export DISPLAY=:0 && /home/idachev/bin/glipper_restart.sh > /dev/null 2>&1

30 05 * * * /home/idachev/storage/develop/vmkernel-op-licensing-integrate_p4/respin_vim_cat_cron.sh >> /tmp/cron_job_respin_vim_cat_op.log 2>&1
30 12 * * * /home/idachev/storage/develop/vmkernel-op-licensing-integrate_p4/respin_vim_cat_cron.sh >> /tmp/cron_job_respin_vim_cat_op.log 2>&1
30 19 * * * /home/idachev/storage/develop/vmkernel-op-licensing-integrate_p4/respin_vim_cat_cron.sh >> /tmp/cron_job_respin_vim_cat_op.log 2>&1
30 00 * * * /home/idachev/storage/develop/vmkernel-op-licensing-integrate_p4/respin_vim_cat_cron.sh >> /tmp/cron_job_respin_vim_cat_op.log 2>&1

